# nfs4/s6/rpcbind/run
echo "Starting rpcbind..."
/sbin/rpcbind -w