# nfs4/s6/mountd/run
echo "Starting mountd..."
/usr/sbin/rpcall --no-udp --no-nfs-version 2 --no-nfs-version 3